# systemd-units/renew-certificates.timer
[Unit]
Description=Tries to renew Let’s Encrypt certificates every month

[Timer]
OnCalendar=monthly

[Install]
WantedBy=multi-user.target
